# systemd/ics-dm-iot-module-rs.timer
[Unit]
Description=Restart ics-dm-iot-module-rs

[Timer]
OnUnitInactiveSec=10min

[Install]
WantedBy=timers.target