# systemd/iot-module-template-rs.timer
[Unit]
Description=Restart iot-module-template-rs

[Timer]
OnUnitInactiveSec=10min

[Install]
WantedBy=timers.target